# support-files/docker/supervisor/app.conf
[program:app]
command=daphne -b 0.0.0.0 -p 8000 weops_lite.asgi:application
directory=/apps
autostart=true
priority=10
autorestart=true
redirect_stderr=true
stdout_logfile=/var/log/app.log