# support-files/docker/supervisor/beat.conf
[program:beat]
command=celery -A weops_lite beat --loglevel=info
directory=/apps
autostart=true
priority=10
autorestart=true
redirect_stderr=true
stdout_logfile=/var/log/celery.log