# support-files/docker/supervisor/celery.conf
[program:celery]
command=celery -A weops_lite worker --loglevel=info
directory=/apps
autostart=true
priority=10
autorestart=true
redirect_stderr=true
stdout_logfile=/var/log/celery.log